# psconfig/perfsonar-psconfig/unibuild-packaging/deb/perfsonar-psconfig-pscheduler.postrm
#!/bin/bash
 # postrm script
#
# see: dh_installdeb(1)

set -e

if [ $1 = 'remove' ]; then
    systemctl disable psconfig-pscheduler-agent.service  || :
    systemctl daemon-reload  || :
fi